# src/debian/postinst
#!/bin/sh
# postinst script for kurenma
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)

		PACKAGE=kurenma
		cmdSystemV=/etc/init.d/${PACKAGE}

		if which systemctl >/dev/null 2>&1; then

			chmod -x $cmdSystemV || true
			mv -f $cmdSystemV ${cmdSystemV}.disabled
			systemctl enable $PACKAGE
			systemctl daemon-reload
			systemctl start ${PACKAGE}.service || true

		else

			update-rc.d $PACKAGE defaults
			invoke-rc.d $PACKAGE start || true

		fi

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
#
