# src/debian/postinst
#!/bin/sh
# postinst script for kurenma
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


PACKAGE=kurenma
cmdSystemV=/etc/init.d/${PACKAGE}


case "$1" in
    configure)

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;

esac

if [ -x "/etc/init.d/kurenma" ]; then

  update-rc.d kurenma defaults 21 19 >/dev/null
  invoke-rc.d kurenma start || exit $?

fi


while [ -d /run/systemd/system ]; do

  deb-systemd-helper enable ${PACKAGE}.service || systemctl enable ${PACKAGE}.service || true
  systemctl daemon-reload || true
  deb-systemd-invoke start ${PACKAGE}.service || systemctl start ${PACKAGE}.service || break

  if [ -f $cmdSystemV ]; then

    chmod -x $cmdSystemV || true
    mv -f $cmdSystemV ${cmdSystemV}.disabled-by-systemd || true

  fi

  break

done

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
#
