# src/kurenma.service
#
# $Id: kurenma.service 47 2016-09-06 11:40:22+04:00 toor $
#
[Unit]
Description=Kurenma Auth Service for Wireguard Network

[Service]
Type=simple
TimeoutSec=8
ExecStart=/usr/sbin/kurenma start
ExecStop=/usr/sbin/kurenma stop

[Install]
WantedBy=multi-user.target

