# src/kurenma.service
#
# $Id: kurenma.service 48 2016-09-06 17:10:27+04:00 toor $
#
[Unit]
Description=Kurenma Auth Service for Wireguard Network

[Service]
Type=forking
TimeoutSec=8
ExecStart=/usr/sbin/kurenma start
ExecStop=/usr/sbin/kurenma stop

[Install]
WantedBy=multi-user.target

