# src/kurenma.service
#
# $Id: kurenma.service 85 2016-09-28 14:31:17+04:00 toor $
#
[Unit]
Description=Kurenma Auth Service for Wireguard Network

[Service]
Type=forking
TimeoutSec=4
EnvironmentFile=/etc/default/kurenma
ExecStart=/usr/sbin/kurenmactl start
ExecStop=/usr/sbin/kurenmactl stop

[Install]
WantedBy=multi-user.target

