# src/kurenma.service
#
# $Id: kurenma.service 33 2016-08-16 13:33:53+04:00 toor $
#
[Unit]
Description=Kurenma Auth Service for Wireguard Network

[Service]
Type=simple
TimeoutSec=8
Environment=_kurenmas_ctl=start
ExecStart=/usr/sbin/kurenmas
ExecStop=/usr/sbin/kurenmas stop

[Install]
WantedBy=multi-user.target

