# src/kurenma.service
#
# $Id: kurenma.service 32 2016-08-16 00:19:38+04:00 toor $
#
## TODO timeout for child waiting
[Unit]
Description=Kurenma Auth Service for Wireguard Network

[Service]
Type=simple
Environment=_kurenmas_ctl=start
ExecStart=/usr/sbin/kurenmas
ExecStop=/usr/sbin/kurenmas stop

[Install]
WantedBy=multi-user.target

